# uwsm/.config/autostart/org.kde.kdeconnect.daemon.desktop
[Desktop Entry]
Type=Application
Name=KDE Connect Daemon
Exec=kdeconnectd
Icon=kdeconnect
Terminal=false
X-GNOME-Autostart-enabled=true
