# systemd-user/.config/systemd/user/aur-clean.timer
[Unit]
Description=Weekly clean of paru AUR cache

[Timer]
OnBootSec=15min
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
