# systemd-user/.config/systemd/user/aur-trim.timer
[Unit]
Description=Weekly trim of paru build cache

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
