# systemd-user/.config/systemd/user/check-orphans.service
[Unit]
Description=Check for orphaned packages

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-orphans

[Install]
WantedBy=multi-user.target
