# systemd-user/.config/systemd/user/check-orphans.service
[Unit]
Description=Check for orphaned packages

[Service]
Type=oneshot
ExecStart=/bin/bash -c '\
orphans=$(paru -Qtdq 2>/dev/null); \
if [[ -n "$orphans" ]]; then \
  count=$(echo "$orphans" | wc -w); \
  notify-send -u critical "⚠️ Orphaned Packages Found ($count)" "$orphans"; \
fi'
