# systemd-user/.config/systemd/user/failed-units-check.service
[Unit]
Description=Check for failed systemd units

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-units

[Install]
WantedBy=multi-user.target
