# systemd-user/.config/systemd/user/failed-units-check.service
[Unit]
Description=Check for failed systemd units

[Service]
Type=oneshot
PassEnvironment=DISPLAY DBUS_SESSION_BUS_ADDRESS
ExecStart=/usr/bin/env bash -c '\
sys_failed=$(systemctl --failed --no-legend); \
user_failed=$(systemctl --user --failed --no-legend); \
msg=""; \
if [[ -n "$sys_failed" ]]; then \
  msg="[System Units]\n$sys_failed"; \
fi; \
if [[ -n "$user_failed" ]]; then \
  [[ -n "$msg" ]] && msg="${msg}\n\n"; \
  msg="${msg}[User Units]\n$user_failed"; \
fi; \
if [[ -n "$msg" ]]; then \
  notify-send -u critical "⚠️ Systemd Failure" "$msg"; \
else \
  notify-send -u low -t 3000 "✅ System Healthy" "All system and user units are running correctly."; \
fi'
