# systemd-user/.config/systemd/user/failed-units-check.timer
[Unit]
Description=Daily systemd failed units check

[Timer]
OnCalendar=daily
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
