# systemd-user/.config/systemd/user/flatpak-update-low-priority.timer
[Unit]
Description=Run Flatpak update daily

[Timer]
OnBootSec=2m
OnUnitActiveSec=1d
Persistent=true
RandomizedDelaySec=15m 

[Install]
WantedBy=timers.target
