# systemd-user/.config/systemd/user/hyprpaper.service
[Unit]
Description=Hyprland Wallpaper Utility
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/hyprpaper
Type=simple
Restart=on-failure

[Install]
WantedBy=graphical-session.target
